# crontab
# Start Monitoring SSH Authorised Keys
@reboot /bin/bash /root/bin/monitor-ssh-keys.sh &

# Sync Webroot (1h To Get Contents From Backup Server Before Start rsync With --delete Flags)
@reboot sleep 1h; /bin/bash /root/bin/sync-webroot.sh &>> /var/log/sync-webroot.sh.log &

# Sync Webroot (1h To Get Contents From Backup Server Before Start rsync With --delete Flags)
@reboot sleep 1h; /bin/bash /root/bin/sync-with-inotify.sh &>> /var/log/sync-with-inotify.log &

# Apply iptable rules on server reboot
@reboot iptables-restore < /root/bin/iptables.rules

# Check MySQL Replication Status
@daily /bin/bash /root/bin/mysqlstatus.sh | mail -s "$(hostname -f) MySQL Replication Status" mitesh.shah@rtcamp.com

#Check HDD Space
@daily	/bin/bash /root/bin/hddspace.sh

# Correct WebRoot Permissions
0	2	*	*	*	/bin/bash /root/bin/fixpermission.sh


# ClamAV Scan
10	2	*	*	*	/bin/bash /root/bin/customscan.sh


# Rsnapshot MySQL Script
20      *       *       *       *       /bin/bash /root/bin/mysqldump.sh

#Run Rsnapshot hourly
00      *       *       *       *       /usr/bin/rsnapshot hourly

#Run Rsnapshot At 02:00 Everyday
40      02      *       *       *       /usr/bin/rsnapshot daily

#Run Rsnapshot At 06:00 Every Sunday 
50      06      *       *       0       /usr/bin/rsnapshot weekly
                                                    
